# templates/units/lbdock-prom.timer
[Unit]
Description=Prometheus targets update timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
